# etc/sv/lawn/run
#!/bin/sh

PATH=/root/bin:$PATH

exec snooze -D/1 -H6,8,11 announce lawn-alternate
