# etc/sv/drip-monthly/run
#!/bin/sh

PATH=/root/bin:$PATH

exec snooze -D/14 -H9 announce drip-monthly
